# systemd/mpd-notify.service
[Unit]
Description=mpd-notify

[Service]
Type=simple
Restart=always
ExecStart=%h/dotfiles/scripts/mpd-notify

[Install]
WantedBy=default.target
